# sys_config/etc/crontab
6 0 * * * "/home/brlo/.acme.sh"/acme.sh --cron --home "/home/brlo/.acme.sh" > /dev/null
